# services/apache/run
#!/usr/bin/env sh
test -e /var/run/apache2/httpd.pid \
  && rm -f /var/run/apache2/httpd.pid
exec /usr/sbin/httpd -D FOREGROUND
